# override-cron/crontab/root
#root crontabs go here
*   * * * *   wget -O /dev/null web/cron.php
*/5 * * * *   /usr/sbin/logrotate /etc/logrotate.conf
