# override-cron/crontab/root
#root crontabs go here
* * * * * wget -O /dev/null web/cron.php
#* * * * * touch /edocs/minute_cron
