# rdc/docker-cron/crontabs
* *   * * *   wget webREDCAP_WEBROOT_PATHcron.php --spider >/dev/null 2>&1
0 * * * *   /usr/sbin/logrotate /etc/logrotate.conf
